# 08_systemd/logmon/logmon.service
[Unit]
Description=logmon service
[Service]
EnvironmentFile=/etc/sysconfig/logmonconf
ExecStart=/vagrant/logmon/logmon.sh $WORD $LOG
